# CPU.service
[Unit]
Description=CPU mining

[Service]
ExecStart=/root/nimiq/CPU/skypool-node-client

[Install]
WantedBy=multi-user.target