# GPU.service
[Unit]
Description=GPU mining

[Service]
ExecStart=/root/nimiq/GPU/sushipool-gpu-miner

[Install]
WantedBy=multi-user.target