# GPU.service
[Unit]
Description=GPU mining

[Service]
ExecStart=/nimiq/GPU/skypool-node-client

[Install]
WantedBy=multi-user.target