# NONCER.service
[Unit]
Description=noncer mining

[Service]
ExecStart=/nimiq/noncer/mine.sh

[Install]
WantedBy=multi-user.target