# NONCER.service
[Unit]
Description=noncer mining

[Service]
ExecStart=/nimiq/noncer/noncerpro

[Install]
WantedBy=multi-user.target